# container/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:nginx]
command=nginx -g "daemon off;"
redirect_stderr=true
stdout_logfile=/tmp/nginx.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10

[program:app]
command=gunicorn --reload -w 4 --access-logfile - -t 0 restapp:app
directory=/svc
redirect_stderr=true
stdout_logfile=/tmp/app.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
